# user_scripts/openrc/user/init.d/battery-notify
#!/sbin/openrc-run

name="Dusky Battery Notifier"
description="Battery notification daemon for laptops"
supervisor=supervise-daemon
command="${HOME}/user_scripts/battery/notify/battery_notify.sh"

depend() {
	need dbus
	after hyprland
}

export DBUS_SESSION_BUS_ADDRESS="unix:path=${XDG_RUNTIME_DIR}/bus"