# user_scripts/openrc/user/init.d/hyprsunset
#!/sbin/openrc-run

name="Hyprsunset"
description="Hyprland blue light filter / night light"
supervisor=supervise-daemon
command="/usr/bin/hyprsunset"

depend() {
	need dbus
	after hyprland
}

export DBUS_SESSION_BUS_ADDRESS="unix:path=${XDG_RUNTIME_DIR}/bus"