# user_scripts/openrc/user/init.d/network-meter
#!/sbin/openrc-run

name="Waybar Network Meter"
description="Waybar network speed monitoring daemon"
supervisor=supervise-daemon
command="${HOME}/user_scripts/waybar/network/network_meter_daemon.sh"

depend() {
	need dbus
	after networkmanager
}

export DBUS_SESSION_BUS_ADDRESS="unix:path=${XDG_RUNTIME_DIR}/bus"