# user_scripts/openrc/user/init.d/pipewire-pulse
#!/sbin/openrc-run

name="PipeWire-Pulse"
description="PipeWire PulseAudio compatibility layer"
supervisor=supervise-daemon
command="/usr/bin/pipewire"
command_args="-c pipewire-pulse.conf"

depend() {
	need pipewire dbus
	after hyprland
}