# user_scripts/openrc/user/init.d/pipewire
#!/sbin/openrc-run

name="PipeWire"
description="Audio/video server for Wayland"
supervisor=supervise-daemon
command="/usr/bin/pipewire"

depend() {
	need dbus
	after hyprland
}