# user_scripts/openrc/user/init.d/wireplumber
#!/sbin/openrc-run

name="WirePlumber"
description="PipeWire session/policy manager"
supervisor=supervise-daemon
command="/usr/bin/wireplumber"

depend() {
	need pipewire dbus
	after hyprland
}