# user_scripts/openrc/init.d/battery-notify
#!/sbin/openrc-run

name="Dusky Battery Notifier"
description="Battery notification daemon"
command="/home/dim/dawn/user_scripts/battery/notify/battery_notify.sh"
command_background="yes"
pidfile="/run/${RC_SVCNAME}.pid"
output_log="/var/log/${RC_SVCNAME}.log"
error_log="/var/log/${RC_SVCNAME}.err"

depend() {
    need localmount
    after acpid
}

start_pre() {
    checkpath --directory --owner root:root --mode 0755 /run 2>/dev/null || true
}
