# user_scripts/openrc/init.d/dusky-control-center
#!/sbin/openrc-run

name="Dusky Control Center"
description="Dusky system control center"
command="/home/dim/duskyRC/user_scripts/dusky_system/control_center/dusky_control_center"
command_background="yes"
pidfile="/run/${RC_SVCNAME}.pid"
output_log="/var/log/${RC_SVCNAME}.log"
error_log="/var/log/${RC_SVCNAME}.err"

depend() {
    need localmount
}

start_pre() {
    checkpath --directory --owner root:root --mode 0755 /run 2>/dev/null || true
}
