# user_scripts/openrc/init.d/hypridle
#!/sbin/openrc-run

name="Hypridle"
description="Hyprland idle daemon"
command="/usr/bin/hypridle"
command_background="yes"
pidfile="/run/${RC_SVCNAME}.pid"

depend() {
    need localmount
    after hyprland
}

start_pre() {
    checkpath --directory --owner root:root --mode 0755 /run 2>/dev/null || true
}
