# user_scripts/openrc/init.d/update-checker
#!/sbin/openrc-run

name="Dusky Update Checker"
description="Checks for system updates periodically"
command="/home/dim/duskyRC/user_scripts/update_dusky/update_checker/update_checker.sh"
command_args="--num"
command_background="yes"
pidfile="/run/${RC_SVCNAME}.pid"
output_log="/var/log/${RC_SVCNAME}.log"
error_log="/var/log/${RC_SVCNAME}.err"

depend() {
    need localmount
}

start_pre() {
    checkpath --directory --owner root:root --mode 0755 /run 2>/dev/null || true
}
